# ansible/roles/openbao/templates/openbao-backup.service.j2
# Managed by the openbao Ansible role.
[Unit]
Description=OpenBao daily backup — leader-guarded JSON dump to backup-server
After=network-online.target openbao.service
Wants=network-online.target

[Service]
Type=oneshot
User=openbao
Group=openbao
ExecStart=/usr/local/sbin/openbao-backup
